# ops/systemd/as6-production-backup.service
[Unit]
Description=AS6 production full backup (PostgreSQL and attachments)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/ai-platform
ExecStart=/var/www/ai-platform/ops/bin/as6-production-backup-full-v1
User=root
Group=root
Nice=10
